# deploy/ansible/roles/services/templates/api.conf
description "daemon for api"

start on startup
stop on shutdown
respawn

env NODE_ENV={{env}}

script
  /usr/bin/npm run start-api --prefix {{site_path}}
end script
